# .noodle/com.noodle.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noodle.scheduler</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mohlt/.openclaw/workspace/.noodle/noodle.sh</string>
        <string>run</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mohlt/.openclaw/workspace/.noodle/logs/scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mohlt/.openclaw/workspace/.noodle/logs/scheduler-error.log</string>
</dict>
</plist>
