# scripts/com.mohltbot.ghost-shift.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mohltbot.ghost-shift</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mohlt/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>/Users/mohlt/.openclaw/workspace/ghost-shift-work/scripts/ghost-shift-v3.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mohlt/.openclaw/workspace</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/mohlt</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <!-- Run at 12:00 PM (noon) AND 12:00 AM (midnight) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/mohlt/.openclaw/workspace/logs/ghost-shift.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mohlt/.openclaw/workspace/logs/ghost-shift-error.log</string>
</dict>
</plist>
