# scripts/com.mohltbot.event-scout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mohltbot.event-scout</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/mohlt/.openclaw/workspace/scripts/event-scout.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mohlt/.openclaw/workspace</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/mohlt</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/mohlt/.openclaw/workspace/logs/event-scout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mohlt/.openclaw/workspace/logs/event-scout-error.log</string>
</dict>
</plist>
