# scripts/openclaw-debugger-shift3.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.debugger.shift3</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mohlt/.openclaw/workspace/scripts/openclaw-debugger-cron.sh</string>
        <string>3</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/mohlt/.openclaw/workspace/logs/debugger-shift3.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mohlt/.openclaw/workspace/logs/debugger-shift3-error.log</string>
</dict>
</plist>
